# com.reed.mitreserve.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.reed.mitreserve</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/reed/workspace/mitreserve/.venv/lib/python3.13/site-packages</string>
    </dict>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/reed/workspace/mitreserve/.venv/bin/python3</string>
        <string>/Users/reed/workspace/mitreserve/mitreserve.py</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardErrorPath</key>
    <string>/tmp/mitreserve.err</string>

    <key>StandardOutPath</key>
    <string>/tmp/mitreserve.out</string> 
</dict>
</plist>
